# config/com.clawoss.pr-ledger-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.clawoss.pr-ledger-sync</string>

    <key>Comment</key>
    <string>ClawOSS PR Ledger Sync — keeps pr-ledger.md in sync with GitHub</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>__PROJECT_DIR__/scripts/pr-ledger-sync-wrapper.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>60</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/pr-ledger-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/pr-ledger-sync.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key>
      <string>__HOME_DIR__</string>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
  </dict>
</plist>
